# workers/crontab
* * * * * /usr/local/bin/node /Users/Port-O-Bucket/Desktop/Development/hackreactor/2014-01-web-historian/workers/htmlfetcher.js
